# watchd.service
[Unit]
Description=Console Watcher Daemon
Documentation=man:watchd(8)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/watchd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
